# qgc-audio-sender.service
[Unit]
Description=QGC Audio Sender - USB mic to RTSP via mediamtx
After=network.target mediamtx.service
Wants=mediamtx.service

[Service]
Type=simple
ExecStart=/home/raspberry/qgc-audio-sender/send-audio.sh
Restart=on-failure
RestartSec=5
Environment=AUDIO_DEVICE=plughw:2,0
Environment=MEDIAMTX_URL=rtsp://localhost:8554/audio
Environment=BITRATE=64000

[Install]
WantedBy=multi-user.target
